# sources/Kainat-kde-settings_1.0.01_all/etc/profile.d/KOS-xdg.sh
export XDG_CONFIG_DIRS="/usr/share/kainat-kde-settings/xdg:/etc/xdg"
