# sources/Kainat-os-disk-manager/DEBIAN/postinst
#!/bin/bash
set -e

# Enable the service after installation
systemctl enable disk-links.service
