# sources/Kainat-os-disk-manager/DEBIAN/postinst
#!/bin/sh
set -e

# Path to the desktop file
SOURCE_FILE="/etc/skel/Desktop/computer.desktop"
SOURCE_FOLDER="/etc/skel/.local"

# Ensure the source file exists
if [ ! -f "$SOURCE_FILE" ]; then
    echo "Source file $SOURCE_FILE does not exist."
    exit 1
fi

# Iterate over each user with UID >= 1000 (standard for regular users)
getent passwd | awk -F: '$3 >= 1000 && $3 < 65534 { print $1 ":" $6 }' | while IFS=: read -r user home; do
    # Skip if home directory doesn't exist
    if [ ! -d "$home" ]; then
        echo "Home directory $home for user $user does not exist. Skipping."
        continue
    fi

    # Create Desktop directory if it doesn't exist
    mkdir -p "$home/Desktop"

    # Copy the desktop file
    cp "$SOURCE_FILE" "$home/Desktop/"
    cp -r "$SOURCE_FOLDER" "$home/"

    # Set ownership
    chown "$user":"$user" "$home/Desktop/computer.desktop"

    echo "Installed computer.desktop to $home/Desktop for user $user."
done
if [  -d "/computer" ]; then
    chmod 777 /computer
fi
update-mime-database /usr/share/mime



# This script is run by dpkg after package installation/configuration.

case "$1" in
  configure)
    # only try to use systemd if it's available
    if command -v systemctl >/dev/null 2>&1; then
      # reload systemd to pick up your unit file
      systemctl daemon-reload || true

      # enable the service so it runs on every boot
      systemctl enable create-computer-dir.service || true

      # (optionally) start it right now
      systemctl start create-computer-dir.service || true
    fi
    ;;
esac

exit 0
