# sources/Kainat-os-disk-manager/etc/systemd/system/disk-links.path.service
[Unit]
Description=Watch for mount table changes to update disk links

[Path]
PathChanged=/proc/self/mountinfo

[Install]
WantedBy=multi-user.target
